# supervisord.ini
[supervisord]
nodaemon = true
loglevel = debug

[program:pypy]
command = /app/vectors/simserver.py
redirect_stderr = true